# ansible/roles/configure_user/files/gpg-agent.service
[Unit]
Description=GPG agent with ssh-agent support

[Service]
ExecStart=/usr/bin/gpg-connect-agent
