# ansible/roles/configure_user/files/link_dotfiles.timer
[Unit]
Description=Run dotfiles/config.py every 15 minutes

[Timer]
OnCalendar=*-*-* *:00/15:00
Persistent=false

[Install]
WantedBy=timers.target