# ansible/roles/configure_user/files/link_dotfiles.service
[Unit]
Description=Run dotfiles/config.py

[Service]
ExecStart=%h/dotfiles/config.py ln
Type=oneshot