# ansible/roles/configure_user/files/mount-nas-sshfs.service
[Unit]
Description=Mount nas with sshfs
After=network-online.target

[Service]
ExecStart=%h/bin/mount-nas-sshfs

[Install]
WantedBy=default.target
